# conf/rbus_heartbeat.timer
[Unit]
Description=rbus heart beat timer

[Timer]
OnBootSec=6min
OnUnitActiveSec=2min

[Install]
WantedBy=timers.target
